# service/yui-bot.service.in
# Copyright (c) 2025 Guppy Girl Genetics Software
# SPDX-License-Identifier: BSD-2-Clause
# See LICENSE file for full text.
#
# Systemd Unit File Template for @PACKAGE_NAME@
# Processed by configure script.

[Unit]
Description=Yui Discord Bot Service (@PACKAGE_NAME@)
After=network.target

[Service]
Type=simple
User=@installuser@
Group=@installgroup@
WorkingDirectory=@FINAL_PKG_DATA_DIR@

# --- Execution ---
ExecStart=@FINAL_PYTHON_EXEC@ @FINAL_PKG_DATA_DIR@/yui_bot.py --config @FINAL_ENV_FILE@ --pidfile @FINAL_PID_FILE@ --log-level INFO
PIDFile=@FINAL_PID_FILE@
Restart=on-failure
RestartSec=5s

# --- Environment ---
Environment=PYTHONUNBUFFERED=1

# --- Security Hardening (One directive per line) ---
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=@FINAL_APPRUN_DIR@
# ProtectKernelTunables=yes
# ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
SystemCallArchitectures=native
# Combine SystemCallFilter rules (adjust if needed)
SystemCallFilter=@system-service ~@privileged ~@resources ~@raw-io ~@reboot ~@swap

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
