# service/yui-bot.service
# Copyright (c) 2025 Guppy Girl Genetics Software
# SPDX-License-Identifier: BSD-2-Clause
# See LICENSE file for full text.
#
# Systemd Unit File Template for yui-bot
# Processed by configure script.

[Unit]
Description=Yui Discord Bot Service (yui-bot)
After=network.target

[Service]
Type=simple
User=yui-bot
Group=yui-bot
WorkingDirectory=${datadir}/${PACKAGE_NAME}

# --- Execution ---
ExecStart=python3 ${datadir}/${PACKAGE_NAME}/yui_bot.py --config ${prefix}/etc/yui-bot/.env --pidfile ${localstatedir}/run/yui-bot/yui-bot.pid --log-level INFO
PIDFile=${localstatedir}/run/yui-bot/yui-bot.pid
Restart=on-failure
RestartSec=5s

# --- Environment ---
Environment=PYTHONUNBUFFERED=1

# --- Security Hardening ---
# AmbientCapabilities=; CapabilityBoundingSet=;
NoNewPrivileges=yes; PrivateTmp=yes; ProtectSystem=strict; ProtectHome=yes;
ReadWritePaths=${localstatedir}/run/yui-bot; #ProtectKernelTunables=yes; ProtectKernelModules=yes;
ProtectControlGroups=yes; RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6;
RestrictNamespaces=yes; RestrictRealtime=yes; SystemCallArchitectures=native;
SystemCallFilter=@system-service; SystemCallFilter=~@privileged @resources @raw-io @reboot @swap;

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
